# docker/base/supervisor/conf.d/php5.conf
[program:php5]
command=/usr/sbin/php5-fpm -F
stdout_logfile=syslog
stderr_logfile=syslog
autostart=true
autorestart=true