# docker/base/supervisor/conf.d/rsyslogd.conf
[program:rsyslogd]
command=/usr/sbin/rsyslogd -n
stdout_logfile=syslog
stderr_logfile=syslog
autostart=true
autorestart=true
priority=1