# docker/supervisor/conf.d/nginx.conf
[program:nginx]
command=/usr/sbin/nginx
stdout_logfile=syslog
stderr_logfile=syslog
autostart=true
autorestart=true