# docker/supervisor/conf.d/php.conf
[program:php]
command=/usr/sbin/php-fpm7.0
stdout_logfile=syslog
stderr_logfile=syslog
autostart=true
autorestart=true