# docker/supervisor/conf.d/redis.conf
[program:redis]
command=/usr/bin/redis-server /etc/redis/redis.conf
user=redis
stdout_logfile=syslog
stderr_logfile=syslog
autostart=false
autorestart=true